# examples/com.localdrop.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS startup for LocalDrop (optional).

  Install:
    cp com.localdrop.plist.example ~/Library/LaunchAgents/com.localdrop.plist
    # edit the paths below, then:
    launchctl load ~/Library/LaunchAgents/com.localdrop.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.localdrop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/YOURNAME/projects/localdrop/main.py</string>
        <string>daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOURNAME/projects/localdrop</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
